# docker/files/nginxstd.conf.j2
[program:nginx]
command = /usr/sbin/nginx 
autostart = true
stdout_logfile = /nginx.log
redirect_stderr = true
autorestart = true 