# docker/files/phpstd.conf.j2
[program:php]
command = /usr/sbin/php-fpm7.4 
autostart = true
stdout_logfile = /phpstd.log
redirect_stderr = true
autorestart = true 